# systemd/ai-radio-station-id.timer
[Unit]
Description=AI Radio Station - Station ID Enqueue Timer
Documentation=https://github.com/yourusername/ai-radio
Requires=ai-radio-station-id.service

[Timer]
# Run 60 seconds early (:14, :29, :44) to enqueue before target times (:15, :30, :45)
# Gives time for generation without interrupting normal queue flow
OnBootSec=3min
OnCalendar=*:14,29,44:00

# Persistent (catches up if system was down)
Persistent=true

# Small randomization to avoid exact clustering
RandomizedDelaySec=5

[Install]
WantedBy=timers.target
