# systemd/ai-radio-enqueue.service
[Unit]
Description=AI Radio Station - Music Enqueue Service
After=network.target ai-radio-liquidsoap.service
Requires=ai-radio-liquidsoap.service

[Service]
Type=oneshot
User=radio
Group=radio
WorkingDirectory=/srv/ai_radio

# Set PYTHONPATH so Python can find ai_radio package
Environment="PYTHONPATH=/srv/ai_radio/src"

ExecStart=/srv/ai_radio/.venv/bin/python /srv/ai_radio/scripts/enqueue_music.py

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-radio-enqueue

# Nice level to avoid starving Liquidsoap
Nice=10

[Install]
WantedBy=multi-user.target
