# systemd/ai-radio-export-nowplaying.timer
[Unit]
Description=AI Radio Now Playing Export Timer (Fallback)
Documentation=https://github.com/your-username/ai-radio-station

[Timer]
# Fallback timer - runs every 2 minutes as safety net
# Normal operation: exports triggered immediately by record_play.py
# This timer only catches cases where triggered export fails
OnBootSec=30s
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target
