# systemd/ai-radio-schedule-station-id.service
[Unit]
Description=AI Radio Station - Station ID Scheduler
After=network.target ai-radio-liquidsoap.service
BindsTo=ai-radio-liquidsoap.service

[Service]
Type=oneshot
User=ai-radio
Group=ai-radio
WorkingDirectory=/srv/ai_radio
ExecStart=/srv/ai_radio/.venv/bin/python /srv/ai_radio/scripts/schedule_station_id.py
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
