# systemd/generate-break.timer
[Unit]
Description=AI Radio Break Generation Timer
Documentation=https://github.com/your-username/ai-radio-station
Requires=generate-break.service

[Timer]
# Run every 50 minutes (matches break_freshness_minutes config)
OnCalendar=*:00,50:00

# Run immediately if system was off during scheduled time
Persistent=true

# Randomize start time by up to 2 minutes to avoid precise timing patterns
RandomizedDelaySec=120

# Start timer immediately on system boot (with random delay)
OnBootSec=5min

[Install]
WantedBy=timers.target
